# moonsig.crontab
# m h  dom mon dow   command
MAILTO=neilk@brevity.org
PROJECTS=...
00 */6 * * * $PROJECTS/moonsig/moonsig neilk@brevity.org   > $HOME/signature.neilk@brevity.org
00 */6 * * * $PROJECTS/moonsig/moonsig neilk@neilk.net     > $HOME/signature.neilk@neilk.net
