# moonsig.crontab
# m h  dom mon dow   command
MAILTO=neilk@brevity.org
HOME=/home/neilk
00 */6 * * * $HOME/bin/sig neilk@brevity.org   > $HOME/signature.neilk@brevity.org
00 */6 * * * $HOME/bin/sig neilk@wikimedia.org > $HOME/signature.neilk@wikimedia.org
